# salt/backup/crontab
0 0 * * * {{ pillar.backup.cron_run_script }}
