# salt/syslog/crontab
* * * * * {{ pillar.syslog.inject_log_service_path }} --dry-run 2>&1 /usr/bin/vlogger -t crond-inject-log-service
